# hyprland/.config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Wayland compositor session
Documentation=man:systemd.special(7)
Requires=hyprland.service
Wants=hyprland-session-xdg-autostart.target
After=graphical-session-pre.target
BindsTo=graphical-session.target
Before=graphical-session.target
PropagatesStopTo=graphical-session.target
Conflicts=hyprland-session-shutdown.target
Before=hyprland-session-shutdown.target
StopWhenUnneeded=yes

[Install]
# WantedBy=default.target
# Comment out the above line to prevent auto-start on login
# Enable manually with: systemctl --user enable hyprland-session.target
# Or start from shell with the start-hyprland-session script
