# systemd/.config/systemd/user/1password.service
[Unit]
Description=1Password Linux desktop client
Documentation=https://support.1password.com/
# Hard link to Waybar: if Waybar restarts, 1Password restarts too.
After=waybar.service
BindsTo=wayland-session.target
PartOf=waybar.service

[Service]
Type=simple
# Fix for Wayland: Force native ozone and clear old DISPLAY hangs
Environment=ELECTRON_OZONE_PLATFORM_HINT=auto

# The "Smart" start: Wait for Waybar to register on DBus before starting
ExecStart=/usr/bin/sh -c 'while ! pgrep -x waybar >/dev/null; do sleep 0.5; done; sleep 1.5; exec /usr/bin/1password --silent'

Restart=on-abnormal
RestartSec=5s
OOMPolicy=continue

[Install]
WantedBy=wayland-session.target
