# systemd/.config/systemd/user/hyprland.service
[Unit]
Description=Hyprland - A dynamic tiling Wayland compositor
Documentation=man:hyprland(1)
BindsTo=wayland-session.target
Before=wayland-session.target graphical-session.target
PropagatesStopTo=wayland-session.target graphical-session.target
After=graphical-session-pre.target
OnSuccess=wayland-session-shutdown.target
OnSuccessJobMode=replace-irreversibly
OnFailure=wayland-session-shutdown.target
OnFailureJobMode=replace-irreversibly
Conflicts=wayland-session-shutdown.target
Before=wayland-session-shutdown.target
CollectMode=inactive-or-failed

[Service]
Type=notify
Environment=LIBVA_DRIVER_NAME=nvidia
Environment=XDG_SESSION_TYPE=wayland
Environment=GBM_BACKEND=nvidia-drm
Environment=__GLX_VENDOR_LIBRARY_NAME=nvidia
ExecStart=/usr/bin/start-hyprland
NotifyAccess=all
# This is critical: it exports Wayland variables so other services (Waybar, etc) can see them
ExecStartPost=/usr/bin/sh -c "systemctl --user import-environment WAYLAND_DISPLAY XDG_CURRENT_DESKTOP HYPRLAND_INSTANCE_SIGNATURE QT_QPA_PLATFORMTHEME; dbus-update-activation-environment --systemd WAYLAND_DISPLAY XDG_CURRENT_DESKTOP HYPRLAND_INSTANCE_SIGNATURE"
Restart=on-failure
KillMode=mixed

[Install]
# This links the compositor to the main session target
WantedBy=wayland-session.target
